# roles/nftables/files/nft-geo-filter.timer
[Unit]
Description=nftables Country Filter Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
