# roles/nftables/files/nft-blacklist.timer
[Unit]
Description=nftables Blacklist Updater Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
