# roles/znc/files/znc.service
[Unit]
Description=ZNC IRC Bouncer

[Service]
User=zncd
ExecStart=/usr/bin/znc -Df

[Install]
WantedBy=multi-user.target
